# pkg/rpm/docker-dnsmasq-updater.service
[Unit]
Description=Poor mans service discovery for docker.
Documentation=https://github.com/svenwltr/docker-dnsmasq-updater
After=docker.service
Requires=docker.service

[Service]
ExecStart=/usr/bin/docker-dnsmasq-updater
